# ansible/roles/execution/templates/gqrl.service.j2
[Unit]
Description=QRL Execution Layer (gqrl)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=qrl
Group=qrl
ExecStart=/usr/local/bin/gqrl \
  --datadir {{ data_dir }}/execution \
  --http --http.addr 0.0.0.0 --http.port {{ gqrl_http_port }} \
  --http.api qrl,net,web3,engine \
  --authrpc.addr 0.0.0.0 --authrpc.port {{ gqrl_auth_port }} \
  --authrpc.jwtsecret {{ jwt_secret_path }} \
  --authrpc.vhosts=* \
  --port {{ gqrl_p2p_port }} \
  --metrics --metrics.addr 0.0.0.0 --metrics.port {{ gqrl_metrics_port }} \
  --networkid {{ chain_id }} \
{% if execution_bootnodes != "" %}
  --bootnodes {{ execution_bootnodes }} \
{% endif %}
  {{ gqrl_extra_flags | default('') }}
Restart=always
RestartSec=5
LimitNOFILE=65536
StandardOutput=append:{{ data_dir }}/logs/gqrl.log
StandardError=append:{{ data_dir }}/logs/gqrl.log

[Install]
WantedBy=multi-user.target