# ansible/roles/spammer/templates/spammer.service.j2
[Unit]
Description=QRL Transaction Spammer (tx-fuzz)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=qrl
Group=qrl
ExecStart=/usr/local/bin/tx-fuzz spam \
  --rpc http://{{ target_execution_ip }}:{{ gqrl_http_port | default(8545) }} \
  --seed {{ spammer_seed }} \
  --txcount {{ spammer_tx_count }} \
  --accounts {{ spammer_accounts }} \
  {{ spammer_extra_flags | default('') }}
Restart=always
RestartSec=5
StandardOutput=append:{{ data_dir }}/logs/tx-fuzz.log
StandardError=append:{{ data_dir }}/logs/tx-fuzz.log

[Install]
WantedBy=multi-user.target